# docker-template/supervisord.conf
[program:mysqld]
command=/usr/bin/mysqld_safe

[program:httpd]
command=/etc/apache2/foreground.sh
stopsignal=6